# templates/etc/profile.d/storm-globus-gridftp.sh.erb
# The contents of this file was generated by Puppet
export GRIDFTP_WITH_DSI="<%=scope['storm::gridftp::gridftp_with_dsi']%>"
export GLOBUS_TCP_PORT_RANGE="<%=scope['storm::gridftp::tcp_port_range_min']%>,<%=scope['storm::gridftp::tcp_port_range_max']%>"