# templates/etc/systemd/system/storm-webdav.service.d/filelimit.conf.erb
[Service]
LimitNOFILE=<%=scope.lookupvar('storm::webdav::storm_limit_nofile')%>