# templates/etc/systemd/system/storm-webdav.service.d/filelimit.conf.erb
[Service]
LimitNOFILE=<%=@storm_limit_nofile%>