# templates/etc/systemd/system/storm-webdav.service.d/storm-webdav.conf.erb
# The contents of this file was generated by Puppet
[Service]
Environment="STORM_WEBDAV_USER=storm"
<%- if scope.lookupvar('storm::webdav::debug') -%>
# To have the JVM debug agent bind on all interfaces use `...,address=*:1044,suspend=n`, otherwise the agent binds only on
# localhost
Environment="STORM_WEBDAV_JVM_OPTS=<%=scope.lookupvar('storm::webdav::jvm_opts')%> -Xdebug -Xrunjdwp:server=y,transport=dt_socket,address=<%=scope.lookupvar('storm::webdav::debug_port')%>,<%- if scope.lookupvar('storm::webdav::debug_suspend') -%>suspend=y<%- else -%>suspend=n<%- end -%>"
<%- else -%>
Environment="STORM_WEBDAV_JVM_OPTS=<%=scope.lookupvar('storm::webdav::jvm_opts')%>"
<%- end -%>
Environment="STORM_WEBDAV_SERVER_ADDRESS=0.0.0.0"
<%- scope.lookupvar('storm::webdav::hostnames').each_with_index do | hostname, idx | -%>
Environment="STORM_WEBDAV_HOSTNAME_<%= idx %>=<%= hostname %>"
<%- end -%>
Environment="STORM_WEBDAV_HTTPS_PORT=<%=scope.lookupvar('storm::webdav::https_port')%>"
Environment="STORM_WEBDAV_HTTP_PORT=<%=scope.lookupvar('storm::webdav::http_port')%>"
Environment="STORM_WEBDAV_CERTIFICATE_PATH=/etc/grid-security/storm-webdav/hostcert.pem"
Environment="STORM_WEBDAV_PRIVATE_KEY_PATH=/etc/grid-security/storm-webdav/hostkey.pem"
Environment="STORM_WEBDAV_TRUST_ANCHORS_DIR=/etc/grid-security/certificates"
Environment="STORM_WEBDAV_TRUST_ANCHORS_REFRESH_INTERVAL=<%=scope.lookupvar('storm::webdav::trust_anchors_refresh_interval')%>"
Environment="STORM_WEBDAV_MAX_CONNECTIONS=<%=scope.lookupvar('storm::webdav::max_concurrent_connections')%>"
Environment="STORM_WEBDAV_MAX_QUEUE_SIZE=<%=scope.lookupvar('storm::webdav::max_queue_size')%>"
Environment="STORM_WEBDAV_CONNECTOR_MAX_IDLE_TIME=<%=scope.lookupvar('storm::webdav::connector_max_idle_time')%>"
Environment="STORM_WEBDAV_SA_CONFIG_DIR=/etc/storm/webdav/sa.d"
Environment="STORM_WEBDAV_JAR=/usr/share/java/storm-webdav/storm-webdav-server.jar"
Environment="STORM_WEBDAV_LOG=/var/log/storm/webdav/storm-webdav-server.log"
Environment="STORM_WEBDAV_OUT=/var/log/storm/webdav/storm-webdav-server.out"
Environment="STORM_WEBDAV_ERR=/var/log/storm/webdav/storm-webdav-server.err"
Environment="STORM_WEBDAV_LOG_CONFIGURATION=/etc/storm/webdav/logback.xml"
Environment="STORM_WEBDAV_ACCESS_LOG_CONFIGURATION=/etc/storm/webdav/logback-access.xml"
Environment="STORM_WEBDAV_VO_MAP_FILES_ENABLE=<%=scope.lookupvar('storm::webdav::vo_map_files_enable')%>"
Environment="STORM_WEBDAV_VO_MAP_FILES_CONFIG_DIR=<%=scope.lookupvar('storm::webdav::vo_map_files_config_dir')%>"
Environment="STORM_WEBDAV_VO_MAP_FILES_REFRESH_INTERVAL=<%=scope.lookupvar('storm::webdav::vo_map_files_refresh_interval')%>"
Environment="STORM_WEBDAV_TPC_MAX_CONNECTIONS=<%=scope.lookupvar('storm::webdav::tpc_max_connections')%>"
Environment="STORM_WEBDAV_TPC_MAX_CONNECTIONS_PER_ROUTE=<%=scope.lookupvar('storm::webdav::tpc_max_connections_per_route')%>"
Environment="STORM_WEBDAV_TPC_VERIFY_CHECKSUM=<%=scope.lookupvar('storm::webdav::tpc_verify_checksum')%>"
Environment="STORM_WEBDAV_TPC_TIMEOUT_IN_SECS=<%=scope.lookupvar('storm::webdav::tpc_timeout_in_secs')%>"
Environment="STORM_WEBDAV_TPC_TLS_PROTOCOL=<%=scope.lookupvar('storm::webdav::tpc_tls_protocol')%>"
Environment="STORM_WEBDAV_TPC_REPORT_DELAY_SECS=<%=scope.lookupvar('storm::webdav::tpc_report_delay_secs')%>"
Environment="STORM_WEBDAV_TPC_ENABLE_TLS_CLIENT_AUTH=<%=scope.lookupvar('storm::webdav::tpc_enable_tls_client_auth')%>"
Environment="STORM_WEBDAV_TPC_PROGRESS_REPORT_THREAD_POOL_SIZE=<%=scope.lookupvar('storm::webdav::tpc_progress_report_thread_pool_size')%>"
Environment="STORM_WEBDAV_AUTHZ_SERVER_ENABLE=<%=scope.lookupvar('storm::webdav::authz_server_enable')%>"
<%- if scope.lookupvar('storm::webdav::authz_server_enable') -%>
Environment="STORM_WEBDAV_AUTHZ_SERVER_ISSUER=<%=scope.lookupvar('storm::webdav::authz_server_issuer')%>"
Environment="STORM_WEBDAV_AUTHZ_SERVER_MAX_TOKEN_LIFETIME_SEC=<%=scope.lookupvar('storm::webdav::authz_server_max_token_lifetime_sec')%>"
Environment="STORM_WEBDAV_AUTHZ_SERVER_SECRET=<%=scope.lookupvar('storm::webdav::authz_server_secret')%>"
<%- end -%>
Environment="STORM_WEBDAV_REQUIRE_CLIENT_CERT=<%=scope.lookupvar('storm::webdav::require_client_cert')%>"
Environment="STORM_WEBDAV_USE_CONSCRYPT=<%=scope.lookupvar('storm::webdav::use_conscrypt')%>"
Environment="STORM_WEBDAV_TPC_USE_CONSCRYPT=<%=scope.lookupvar('storm::webdav::tpc_use_conscrypt')%>"
Environment="STORM_WEBDAV_ENABLE_HTTP2=<%=scope.lookupvar('storm::webdav::enable_http2')%>"