# templates/etc/systemd/system/storm-backend-server.service.d/filelimit.conf.erb
[Service]
LimitNOFILE=<%= scope.lookupvar('storm::backend::storm_limit_nofile') %>