# templates/etc/systemd/system/storm-backend-server.service.d/storm-backend-server.conf.erb
[Service]
# options for the JVM running the StoRM BE server
<%- if scope.lookupvar('storm::backend::debug') -%>
Environment="STORM_BE_JVM_OPTS=<%= scope.lookupvar('storm::backend::jvm_options') %> -Xdebug -Xrunjdwp:server=y,transport=dt_socket,address=<%= scope.lookupvar('storm::backend::debug_port') %>,<%- if scope.lookupvar('storm::backend::debug_suspend') -%>suspend=y<%- else -%>suspend=n<%- end -%>"
<%- else -%>
Environment="STORM_BE_JVM_OPTS=<%= scope.lookupvar('storm::backend::jvm_options') %>"
<%- end -%>
# LCMAPS config file
Environment="LCMAPS_DB_FILE=<%= scope.lookupvar('storm::backend::lcmaps_db_file') %>"
# LCMAPS user mapping policy (defined in $LCMAPS_DB_FILE)
Environment="LCMAPS_POLICY_NAME=<%= scope.lookupvar('storm::backend::lcmaps_policy_name') %>"
# LCMAPS log file
Environment="LCMAPS_LOG_FILE=<%= scope.lookupvar('storm::backend::lcmaps_log_file') %>"
# LCMAPS log verbosity: 0 minimum (default), 5 maximum
Environment="LCMAPS_DEBUG_LEVEL=<%= scope.lookupvar('storm::backend::lcmaps_debug_level') %>"