# templates/etc/systemd/system/storm-backend-server.service.d/storm-backend-server.conf.erb
[Service]
Environment="STORM_BE_JVM_OPTS=<%= scope.lookupvar('storm::backend::jvm_options') %>"
Environment="LCMAPS_DB_FILE=<%= scope.lookupvar('storm::backend::lcmaps_db_file') %>"
Environment="LCMAPS_POLICY_NAME=<%= scope.lookupvar('storm::backend::lcmaps_policy_name') %>"
Environment="LCMAPS_LOG_FILE=<%= scope.lookupvar('storm::backend::lcmaps_log_file') %>"
Environment="LCMAPS_DEBUG_LEVEL=<%= scope.lookupvar('storm::backend::lcmaps_debug_level') %>"