# templates/etc/systemd/system/mariadb.service.d/limits.conf.erb
[Service]
LimitNOFILE=<%=@limit_no_file%>