# templates/etc/systemd/system/mariadb.service.d/limits.conf.erb
[Service]
LimitNOFILE=<%=scope.lookupvar('storm::db::limit_no_file')%>